# development/supervisor.d/api.anlibreeders.com.ini
[program: api.anlibreeders.com]
directory = /home/user/www/anlibreeders-back-end
command = /home/user/www/anlibreeders-back-end/venv/bin/gunicorn --worker-class sync --workers=9 --reload app:app -b localhost:3001 --timeout 600 --log-level debug
autostart = true
autorestart = true
stopasgroup = true
killasgroup = true
user = user
stderr_logfile = /var/log/api.anlibreeders.com/api.anlibreeders.com.err.log
stderr_logfile_maxbytes = 1MB
stderr_logfile_backups = 10
stderr_capture_maxbytes = 1MB
stdout_logfile = /var/log/api.anlibreeders.com/api.anlibreeders.com.out.log
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 10
stdout_capture_maxbytes = 1MB